# config/dwm-panel.service
[Unit]
Description=DWM taskbar refresh deamon

[Install]
WantedBy=default.target

[Service]
Type=simple
#PIDFile=/var/run/dwm-panel.pid
#EnvironmentFile=-/etc/sysconfig/chronyd
ExecStart=/usr/local/bin/dwm-paneld

