# config/notificationd.service
[Unit]
Description=Notification daemon

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
ExecStart=/usr/libexec/notification-daemon

